# fletcher/systemd/fletcher.service.in
[Unit]
Description=Expected Path Solver
After=syslog.target

[Service]
ExecStart=@BIN_INSTALL_DIR@/@BIN_FILE_NAME@ 80
Restart=on-abort

[Install]
WantedBy=multi-user.target
